# minitower_oled.service
[Unit]
Description=Minitower oled Service

[Service]
Type=forking
User=root
# ExecStart=/bin/bash -c 'python3 /usr/local/luma.examples/examples/animated_gif.py &'
ExecStart=/bin/bash -c 'python3 /usr/local/luma.examples/examples/sys_info.py &'
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
